# sql-exporter/prometheus-sql-exporter.path
[Unit]
Description=Trigger prometheus-sql-exporter-restart once the configuration file is available

[Path]
PathExists=/var/run/postgresql/prometheus-sql-exporter.yml

[Install]
WantedBy=multi-user.target
