# sql-exporter/prometheus-sql-exporter-restart.path
[Unit]
Description=Trigger prometheus-sql-exporter-restart once a cluster becomes available

[Path]
PathExistsGlob=/etc/postgresql/*/*/postgresql.conf

[Install]
WantedBy=multi-user.target
