# debian/elephant-shed-omnidb.postinst
#!/bin/sh

set -eu

case $1 in
  configure)
    echo Starting omnidb-server once to initialize /var/lib/omnidb/.omnidb/omnidb-server/ ...
    systemctl start omnidb-server
    sleep 5
    systemctl stop omnidb-server
    systemctl disable omnidb-server
    ;;
esac

#DEBHELPER#
