# debian/elephant-shed-prometheus.postinst
#!/bin/sh

set -e

case $1 in
  configure)
    systemctl daemon-reload

    deb-systemd-helper enable prometheus.service
    deb-systemd-invoke restart prometheus
    ;;
esac

#DEBHELPER#
