# debian/elephant-shed-prometheus-sql-exporter.postinst
#!/bin/sh

set -e

case $1 in
  configure)
    systemctl daemon-reload

    deb-systemd-helper enable prometheus-sql-exporter.service
    deb-systemd-invoke restart prometheus-sql-exporter
    ;;
esac

#DEBHELPER#
