# debian/elephant-shed-omnidb.postinst
#!/bin/sh

set -eu

case $1 in
  configure)
    systemctl disable omnidb-server
    invoke-rc.d omnidb-server stop
    ;;
esac

#DEBHELPER#
