# debian/elephant-shed-omnidb.postinst
#!/bin/sh

set -eu

case $1 in
  configure)
    systemctl enable omnidb-server
    invoke-rc.d omnidb-server restart
    ;;
esac

#DEBHELPER#
